# 4vps/etc/supervisor/conf.d/frps.conf
[program:frps]
command=frps -c /etc/frps.ini
autostart=true
autorestart=true
