# orangepi/etc/supervisor/conf.d/frpc.conf
[program:frpc]
command=/usr/local/bin/frpc -c /etc/frpc.ini
autorestart=true

