# orangepi/etc/supervisor/conf.d/opiservices.conf
[program:seafile]
user=seafile
command=/usr/bin/seafile
autorestart=true


[program:cron]
user=root
command=/usr/local/bin/cron.py
autorestart=true

