# orangepi/supervisor/conf.d/dropcache.conf
[program:dropcache]
command=/usr/bin/dropcache
autorestart=true

[program:wifi]
command=/usr/bin/wifi
autorestart=true

