# orangepi/supervisor/conf.d/opiservices.conf
[program:seafile]
user=seafile
command=/usr/bin/seafile
autorestart=true

[program:xunlei]
user=xunlei
command=/usr/bin/xunlei
autorestart=true

